# tax_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-05-01T19:00:00</StartBoundary>
      <ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\USER\AppData\Local\Programs\Python\Python313\python.exe</Command>
      <Arguments>"C:\Users\USER\globes-summary\tax_briefing.py"</Arguments>
      <WorkingDirectory>C:\Users\USER\globes-summary</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
